# fastnet2n2k.service
# systemd unit for the fastnet2n2k bridge.
#
# Install into a dedicated system venv (NOT a pipx ~/.local install — a root-run
# service can't rely on a user's PATH):
#   sudo python3 -m venv /opt/fastnet2n2k
#   sudo /opt/fastnet2n2k/bin/pip install fastnet2n2k
#
# Then edit the --serial device below to match your install, and:
#   sudo cp fastnet2n2k.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now fastnet2n2k.service
#   journalctl -u fastnet2n2k.service -f

[Unit]
Description=fastnet2n2k Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fastnet2n2k

# Bring can0 up at the NMEA2000 bitrate before starting.
# restart-ms 100 lets the controller auto-recover from a bus-off.
ExecStartPre=/bin/sh -c 'ip link set can0 down 2>/dev/null; ip link set can0 up type can bitrate 250000 restart-ms 100'

ExecStart=/opt/fastnet2n2k/bin/fastnet2n2k --serial /dev/ttyUSB0 --channel can0
Restart=always
RestartSec=10

# === RESOURCE LIMITS ===
OOMScoreAdjust=-700
OOMPolicy=continue
MemoryMax=128M
MemoryHigh=96M
TimeoutStopSec=30

# === LOGGING ===
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fastnet2n2k
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
