# fastnet2n2k.service
# systemd unit for the fastnet2n2k bridge.
#
# Install globally with pipx (--global lands the command in /usr/local/bin, where
# a root-run service can reach it — a plain `pipx install` goes to ~/.local/bin,
# which it can't). Needs pipx >= 1.5; run `which fastnet2n2k` to confirm the path.
#   sudo apt install pipx                        # once, if you don't have it
#   sudo pipx install --global fastnet2n2k       # -> /usr/local/bin/fastnet2n2k
#
# Then edit the --serial device below to match your install, and:
#   sudo cp fastnet2n2k.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now fastnet2n2k.service
#   journalctl -u fastnet2n2k.service -f
#
# Update later with: sudo pipx upgrade --global fastnet2n2k

[Unit]
Description=fastnet2n2k Service
After=network.target

[Service]
Type=simple
User=root

# Bring can0 up at the NMEA2000 bitrate before starting.
# restart-ms 100 lets the controller auto-recover from a bus-off.
ExecStartPre=/bin/sh -c 'ip link set can0 down 2>/dev/null; ip link set can0 up type can bitrate 250000 restart-ms 100'

ExecStart=/usr/local/bin/fastnet2n2k --serial /dev/ttyUSB0 --channel can0
Restart=always
RestartSec=10

# === RESOURCE LIMITS ===
OOMScoreAdjust=-700
OOMPolicy=continue
MemoryMax=128M
MemoryHigh=96M
TimeoutStopSec=30

# === LOGGING ===
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fastnet2n2k
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
